# dot_config/systemd/user/gnome-keyring.service
[Unit]
Description=GNOME Keyring Daemon (Secrets)
Documentation=man:gnome-keyring-daemon(1)
# Ensure it starts after the graphical session is ready
After=graphical-session.target
Before=nextcloud-client.service # Ensure it starts before Nextcloud

[Service]
Type=exec
ExecStart=/usr/bin/gnome-keyring-daemon --start --components=secrets
# It's important to set the environment so other applications can find the keyring
Environment=DBUS_SESSION_BUS_ADDRESS=%D/bus
# Keep the service running even if the daemon detaches (Type=exec handles this, but good to be explicit)
RemainAfterExit=yes

[Install]
WantedBy=graphical-session.target
